# templates/networking.service.j2
[Unit]
After=openvswitch-nonetwork.service
After=openvswitch-switch.service
